# systemd/csye6225.service
[Unit]
Description=CSYE6225 Web Application
# 只需要在網路啟動後才啟動
After=network.target
Wants=network-online.target

[Service]
# simple:主程序就是服務本身
Type=simple

# 使用非特權使用者執行
User=csye6225
Group=csye6225

WorkingDirectory=/opt/csye6225

# 啟動應用程式的指令
ExecStart=/usr/bin/java -jar /opt/csye6225/webapp.jar

# 日誌輸出
StandardOutput=journal
StandardError=journal
SyslogIdentifier=csye6225-webapp

# 自動重啟設定
Restart=on-failure
RestartSec=10

# 安全性設定
NoNewPrivileges=true
PrivateTmp=true

[Install]
# 在 multi-user.target 時啟動（相當於開機時）
WantedBy=multi-user.target



